# build/deploy-app/supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command=php-fpm

[program:entrypoint]
command=/bin/bash /entrypoint.sh

[program:mail]
command=/etc/init.d/sendmail start

[program:laravel-queue]
command=php artisan queue:listen --tries=2 --sleep=3
directory=/var/www/html/frontend
stdout_logfile=/var/www/html/frontend/storage/logs/supervisor.log
stderr_logfile=/var/www/html/frontend/storage/logs/supervisor.err.log
autostart=true
autorestart=true

[program:socket-server]
command=nodejs /var/www/html/frontend/socket.js
autostart=true
autorestart=true
stdout_logfile=/var/www/html/frontend/storage/logs/socket-server.log
stderr_logfile=/var/www/html/frontend/storage/logs/socket-server.err.log
