# build/test-deployment-app/supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command=php-fpm

[program:redis]
command=/usr/bin/redis-server

[program:entrypoint]
command=/bin/bash /entrypoint.sh

[program:mail]
command=/etc/init.d/sendmail start

[program:nginx]
command=/etc/init.d/nginx start
