# build/deploy-nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:entrypoint]
command=/bin/bash /entrypoint.sh

